# scripts/systemd/lottery-eod.service
[Unit]
Description=Lottery/Hype bot SAME-DAY EXIT (force-close all positions at ~15:55 ET)
After=network-online.target
Wants=network-online.target

[Service]
# Short-lived: pulls, force-closes every tracked lottery position opened today (cancel trail +
# market close), exits. Pairs with lottery-bot.service (09:44 ET entries). Same-day hold only
# (2026-06-30): avoids the overnight gap bleed. Idle exit 0 if nothing is open / market closed.
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading
ExecStart=/home/ubuntu/trading/scripts/launch_lottery_eod.sh
Restart=no
KillSignal=SIGINT
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
