# scripts/systemd/status.service
[Unit]
Description=ORB status web page (always-on liveness + trade-state dashboard)
# Needs Alpaca reachable to gather trade state; tolerates transient drops (each
# page request re-queries and surfaces errors inline rather than crashing).
After=network-online.target
Wants=network-online.target

[Service]
# Long-running daemon, independent of the timer-fired orb/dualmom runners. It
# reads logs/heartbeat.json (written by the runner) and queries Alpaca; it
# never trades. Safe to restart at any time.
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading
# Public bind (0.0.0.0): the page is intentionally reachable from any device
# with no password. This is acceptable ONLY because it's a PAPER account (the
# equity/positions shown are not real money) and the server is strictly
# read-only (cannot place or cancel orders). Binding alone does NOT expose it --
# reachability also requires the VM firewall (iptables/ufw) AND the Oracle VCN
# security list to allow TCP 8787. If this VM ever runs REAL money, revert to
# 127.0.0.1 + SSH tunnel (or put auth + HTTPS in front).
Environment=STATUS_BIND=0.0.0.0
Environment=STATUS_PORT=8787
# .env (Alpaca keys) is sourced by the script itself via load_env(), matching
# orb.service -- single source of truth, no EnvironmentFile here.
ExecStart=/home/ubuntu/trading/.venv/bin/python /home/ubuntu/trading/live/status_server.py
# Always bring it back: a dashboard that silently died would defeat its own
# purpose. It holds no trading state, so restarting is harmless.
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
