# scripts/systemd/tsmom-rebalance.timer
[Unit]
Description=Fire the monthly TSMOM rebalance on the 2nd-4th of each month (10:15 ET, during RTH)

[Timer]
# Monthly rebalance. Fires the 2nd/3rd/4th at 10:15 ET -- the first trading day rebalances, the others
# are idempotent no-ops (covers weekends/holidays landing on the 2nd). Persistent = catch up if the VM
# was down at the fire time.
OnCalendar=Mon-Fri *-*-2,3,4,5,6 10:15 America/New_York
Persistent=true
Unit=tsmom-rebalance.service

[Install]
WantedBy=timers.target
